# root/etc/s6-overlay/s6-rc.d/p01-file2env/up
/etc/s6-overlay/s6-rc.d/p01-file2env/run
